# .config/autostart/CAMERA_MQTT_CLIENT.desktop
[Desktop Entry]
Name=CAMERA_CLIENT
Exec=lxterminal -e sudo python3 /home/pi/camera_client/camera_client.py
Type=Application
